# scripts/cerulean.service
[Unit]
Description=Cerulean Next.js
After=network.target

[Service]
Type=simple
User=cerulean
Group=cerulean
WorkingDirectory=/home/cerulean/cerulean
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# npm, yerel node_modules/.bin/next'i kullanır (dist yolu sürüme göre değişebilir)
ExecStart=/usr/bin/npm run start -- -H 0.0.0.0 -p 3000
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
