# scripts/cerulean.service
[Unit]
Description=Cerulean Next.js
After=network.target

[Service]
Type=simple
User=cerulean
Group=cerulean
WorkingDirectory=/home/cerulean/cerulean
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
# Doğrudan next binary — npm sarmalayıcısı bazen hemen çıkıyor gibi görünür; curl erken 000 verebilir
ExecStart=/usr/bin/node /home/cerulean/cerulean/node_modules/next/dist/bin/next start -H 0.0.0.0 -p 3000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
